# redhawk-cluster-development/sshd.ini
[program:sshd]
command=/usr/sbin/sshd -D
user=root
priority=3